# apps/public-search-bot/deploy/public-search-bot.service.example
[Unit]
Description=InfinityLinks Public Search Bot
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/infinitylinks-public-search-bot
EnvironmentFile=/opt/infinitylinks-public-search-bot/.env
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=5
# Ensure /opt/infinitylinks-public-search-bot/data is writable by this user.
User=infinitylinks
Group=infinitylinks
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/infinitylinks-public-search-bot/data

[Install]
WantedBy=multi-user.target
