# apps/public-search-bot/deploy/public-search-bot.service.example
[Unit]
Description=InfinityLinks Public Search Bot
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/publicinfinity
EnvironmentFile=/etc/infinitylinks/public-search-bot.env
Environment=NODE_OPTIONS=--dns-result-order=ipv4first
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=5
# Ensure /var/lib/infinitylinks is writable by this user.
User=infinitylinks
Group=infinitylinks
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/infinitylinks

[Install]
WantedBy=multi-user.target
